# config/moacc_ifs12/var/spool/cron/crontabs/nort
# Nort's crontab
*/5 * * * * /home/citco2/bin/Email/doreport
